# supervisord.conf
[program:views]
command: /usr/bin/python -m views
directory: /home/ubuntu/chatbot
autostart: true
stdout_logfile: /tmp/sdbuglog