# supervisord.conf
[program:bookapp]
command: /usr/bin/python -m bookapp
directory: /home/ubuntu/bookapp
autostart: true